# app_bin/ebfdemo/DEBIAN/postinst
#!/bin/bash
systemctl enable ebf-qtdemo
